# deployment/eufmd-nexus-api.service
[Unit]
Description=EuFMD Nexus API
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/var/www/eufmd-nexus/backend
Environment=PATH=/usr/local/bin:/usr/bin:/bin
EnvironmentFile=/etc/systemd/system/eufmd-nexus-api.env
ExecStart=uvicorn main:app --host 0.0.0.0 --port 5800
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target